# dockerfile/jnlp-slave/supervisord.conf
[supervisord]
nodaemon=true

[program:jenkins-slave]
command=/usr/local/bin/jenkins-slave

[program:sshd]
command=/usr/sbin/sshd -D
